# server/stuff-tracker.service
[Unit]
Description=Stuff Tracker Backend Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/app
ExecStart=/home/ubuntu/app/server
Restart=always
RestartSec=3
Environment=RUST_LOG=info
Environment=DATABASE_URL=sqlite:data.db
Environment=PORT=3000
Environment=HOST=127.0.0.1
# Ensure uploads directory exists and is writable
ExecStartPre=/bin/mkdir -p /home/ubuntu/app/uploads

[Install]
WantedBy=multi-user.target
